# .supervisord/api.conf
[program:api]
command=start-kernelgateway.sh
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
autostart=true
autorestart=true
exitcodes=0
